# deploy/combofinder-api.service
[Unit]
Description=ComboFinder API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/combofinder

Environment=PORT=3001
Environment=NODE_ENV=production
Environment=DATABASE_URL=postgresql://combofinder:SIFIRDAN_GUCLU_SIFRE@localhost:5432/combofinder
Environment=SESSION_SECRET=en-az-32-karakterli-rastgele-bir-string-yazin
Environment=ADMIN_USERNAME=admin
Environment=ADMIN_PASSWORD=guclu-sifrenizi-yazin

ExecStart=/usr/bin/node /var/www/combofinder/artifacts/api-server/dist/index.mjs
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
